# crontab
0 * * * * /Users/camillebaldock/Code/Quantified\ Camille/dashboard/export_task.sh
